# cron/wordle.cron
# Sync used Wordle answers once per day at 04:00 local time
0 4 * * * python /app/cli.py used-sync >> /app/data/cron.log 2>&1